# systemd/system/motdgen.service
[Unit]
Description=Generate /run/flatcar/motd
Before=systemd-user-sessions.service

[Service]
Type=oneshot
ExecStart=/usr/lib/flatcar/motdgen
